# slacker.timer
[Unit]
Description=Run foo weekly and on boot

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min 
WakeSystem=true

[Install]
WantedBy=timers.target
